# scripts/phone-server/units/idle-autostop.service
[Unit]
Description=Stop instance when idle
[Service]
Type=oneshot
ExecStart=/usr/local/bin/idle-autostop.sh
